# backend/docker/prod/supervisord/supervisord-80-pulse.conf
[program:pulse-check]
command=/usr/local/bin/php -d variables_order=EGPCS artisan pulse:check
environment=PULSE_SERVER_NAME="hermes"
directory=/srv/www
user=aurora
group=aurora
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:pulse-work]
command=/usr/local/bin/php artisan pulse:work
directory=/srv/www
user=aurora
group=aurora
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
